# monitoring/systemd/prometheus.service
# Prometheus Service for AI Backend Unified Infrastructure
#
# Metrics collection and storage for:
# - LiteLLM gateway performance
# - Provider health and performance
# - System resource utilization
# - Redis cache metrics
#
# Installation:
#   mkdir -p ~/.config/systemd/user
#   cp monitoring/systemd/prometheus.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable prometheus.service
#   systemctl --user start prometheus.service
#
# Logs:
#   journalctl --user -u prometheus.service -f

[Unit]
Description=Prometheus Monitoring Server
Documentation=https://prometheus.io/docs/
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/LAB/ai/backend/ai-backend-unified/monitoring/prometheus

# Prometheus binary (assumes installed via package manager or downloaded)
# Update path if installed elsewhere
ExecStart=/usr/local/bin/prometheus \
  --config.file=%h/LAB/ai/backend/ai-backend-unified/monitoring/prometheus/prometheus.yml \
  --storage.tsdb.path=%h/LAB/ai/backend/ai-backend-unified/monitoring/prometheus/data \
  --storage.tsdb.retention.time=15d \
  --storage.tsdb.retention.size=10GB \
  --web.listen-address=127.0.0.1:9090 \
  --web.console.templates=/usr/share/prometheus/consoles \
  --web.console.libraries=/usr/share/prometheus/console_libraries \
  --web.enable-lifecycle

# Restart policy
Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/LAB/ai/backend/ai-backend-unified/monitoring/prometheus/data

# Resource limits
MemoryLimit=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=prometheus

[Install]
WantedBy=default.target
